# dronestrike-backend.service
[Unit]
Description=DroneStrike v2 Django Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/dronestrike-v2/django_backend
Environment=PATH=/var/www/dronestrike-v2/venv/bin
ExecStart=/var/www/dronestrike-v2/venv/bin/python manage.py runserver 0.0.0.0:8000 --noreload
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target